# components/FIFA22/airootfs.x86_64/etc/systemd/system/onion-greeter.service
[Unit]
Description=Tor control port filter proxy
Documentation=https://tails.boum.org/contribute/design/

[Service]
Type=simple
ExecStart=python /usr/lib/onion-greeter

[Install]
WantedBy=multi-user.target

